# src/nvpn-deamon.service
[Unit]
Description=The NordVPN Deamon
After=network.target

[Service]
ExecStart = /opt/nvpn/bin/vpn_deamon.py

[Install]
WantedBy=multi-user.target